# vm/myapp.service
[Unit]
Description=Gunicorn instance para servir qualitas-fraudes-data
After=network.target

[Service]
# El servicio se ejecuta como usuario de sistema
User=www-data
Group=appusers

# Rutas actualizadas
WorkingDirectory=/opt/qualitas-fraudes-data
Environment="PATH=/opt/qualitas-fraudes-data/venv/bin"

# Configuración optimizada de Gunicorn
ExecStart=/opt/qualitas-fraudes-data/venv/bin/gunicorn \
    --workers=8 \
    --threads=2 \
    --worker-class=gthread \
    --worker-connections=1000 \
    --timeout=120 \
    --max-requests=1000 \
    --max-requests-jitter=50 \
    -b 127.0.0.1:8000 \
    wsgi:app

# Límites de recursos
LimitNOFILE=65536
TimeoutStopSec=20
# Ajusta estos valores según la memoria y CPU disponibles en tu servidor
MemoryLimit=4G
CPUQuota=400%

# Reinicio automático
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target